# websockets.conf
[program:websockets]
command=php /app/artisan websocket:serve --port=6020
numprocs=1
autostart=true
autorestart=true
user=root