# websockets.conf
[program:websockets]
command=php /app/artisan websocket:serve
numprocs=1
autostart=true
autorestart=true
user=root